# service/postgres/run
#!/bin/sh

exec 2>&1
install -d -m 0750 -o "$USER" -g "$GROUP" "$PGDATA"
install -d -m 1775 -o "$USER" -g "$GROUP" /run/postgresql /var/log/postgresql
install -d -m 0750 -o "$USER" -g "$GROUP" "$HOME/pg_arclog"
rm -f "$PGDATA/postmaster.pid"
realpath "$0"
set -ex
test -d "$PGDATA/base" || chpst -u "$USER":"$GROUP" /etc/service/postgres/init
test -d "$PGDATA/base"
exec chpst -u "$USER":"$GROUP" postmaster
